# core-container/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777



[supervisord]
nodaemon=false
pidfile=/tmp/supervisor.pid
logfile = /dev/stdout
logfile_maxbytes = 0

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket



[program:airflow_webserver]
command=airflow webserver -p 8888
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/admin",USER="admin"
startsecs=60
autorestart=unexpected


[program:airflow_scheduler]
command=/bin/bash -c "airflow scheduler"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/admin",USER="admin"
startsecs=60
autorestart=unexpected


[program:mlflow_server]
command=/bin/bash -c "mlflow server --port 4040 --static-prefix /user/admin/mlflow --backend-store-uri postgresql://airflow:airflow@postgres:5432/airflow --default-artifact-root /home/admin/mlflow"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/admin",USER="admin"
startsecs=10
autorestart=unexpected